# deploy/modastack.service
[Unit]
Description=Modastack manager
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/dev/modastack
ExecStart=%h/dev/modastack/.venv/bin/modastack start
Restart=always
RestartSec=5
Environment=PATH=%h/dev/modastack/.venv/bin:%h/.bun/bin:%h/.npm-global/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=default.target
